# systemd/dci-rdo-feeder.socket
[Unit]
Description=Start RDO update on demand

[Socket]
ListenStream=80

[Install]
WantedBy=sockets.target
